# .config/systemd/system/snapshot.service
[Unit]
Description=Create snapshot of the system

[Service]
ExecStart=/home/saverio/.local/bin/snapshot.sh
Type=oneshot
